# confluent_server/systemd/confluent.service
# IBM(c) 2015 Apache 2.0
[Unit]
Description=Confluent hardware manager 

[Service]
Type=forking
PIDFile=/var/run/confluent/pid
ExecStart=/opt/confluent/bin/confluent
ExecStop=/opt/confluent/bin/confetty shutdown /

[Install]
WantedBy=multi-user.target

