# confluent_osdeploy/el7-diskless/profiles/default/scripts/firstboot.service
[Unit]
Description=First Boot Process
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=/opt/confluent/bin/firstboot.sh

[Install]
WantedBy=multi-user.target

